# config/tacquito.service
[Unit]
Description=Tacquito TACACS+ Server
Documentation=https://github.com/facebookincubator/tacquito
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tacquito
Group=tacquito
# Defaults; override via /etc/systemd/system/tacquito.service.d/tacctl-overrides.conf
# (managed by `tacctl config listen` and `tacctl config loglevel`).
Environment="TACQUITO_NETWORK=tcp"
Environment="TACQUITO_ADDRESS=:49"
Environment="TACQUITO_LEVEL=20"
ExecStart=/usr/local/bin/tacquito \
    -config /etc/tacquito/tacquito.yaml \
    -network ${TACQUITO_NETWORK} -address ${TACQUITO_ADDRESS} \
    -acct-log-path /var/log/tacquito/accounting.log \
    -level ${TACQUITO_LEVEL}
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/tacquito
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes

# Port 49 requires CAP_NET_BIND_SERVICE since it is < 1024
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
